# scripts/etc/systemd/system/qmanager-cfun-fix.service
# /lib/systemd/system/qmanager-cfun-fix.service
[Unit]
Description=QManager CFUN=1 Boot Fix (Radio Recovery)
After=qmanager-setup.service
Before=qmanager-poller.service

[Service]
Type=oneshot
ExecStartPre=/bin/sh -c 'i=0; while [ $i -lt 30 ] && [ ! -e /dev/smd11 ]; do sleep 1; i=$((i+1)); done; [ -e /dev/smd11 ] || exit 1'
ExecStart=/usr/bin/qmanager_cfun_fix
RemainAfterExit=yes
TimeoutStartSec=45

[Install]
WantedBy=multi-user.target
