# scripts/etc/systemd/system/qmanager-traffic.service
# /lib/systemd/system/qmanager-traffic.service
[Unit]
Description=QManager Live Traffic Counter Daemon
After=network.target qmanager-setup.service

[Service]
Type=simple
ExecStart=/usr/bin/qmanager_traffic
EnvironmentFile=-/etc/qmanager/environment
TimeoutStopSec=10
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=3600
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
