# scripts/etc/systemd/system/qmanager-ping.service
# /lib/systemd/system/qmanager-ping.service
[Unit]
Description=QManager Ping Daemon (Rust)
After=network.target qmanager-setup.service

[Service]
Type=simple
ExecStart=/usr/bin/qmanager_ping
# Documented relaxed-profile defaults — preserve today's 5s probe behavior on
# fresh starts before /etc/qmanager/ping_profile.json is read. Operator
# overrides via EnvironmentFile= below win because systemd processes these
# directives in declaration order and later assignments override earlier ones
# for the same variable.
Environment=PING_PROFILE=relaxed
Environment=PING_INTERVAL=5
Environment=FAIL_SECS=15
Environment=RECOVER_SECS=10
Environment=INTERCEPT_SECS=8
Environment=HISTORY_SECS=300
Environment=PING_TARGET_1=http://www.gstatic.com/generate_204
Environment=PING_TARGET_2=http://cp.cloudflare.com/
Environment=CARRIER_FILE=/sys/class/net/rmnet_data0/carrier
EnvironmentFile=-/etc/qmanager/environment
TimeoutStopSec=10
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=3600
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
