# scripts/etc/systemd/system/qmanager-ping.service
# /lib/systemd/system/qmanager-ping.service
[Unit]
Description=QManager Ping Daemon
After=network.target qmanager-setup.service

[Service]
Type=simple
ExecStart=/usr/bin/qmanager_ping
EnvironmentFile=-/etc/qmanager/environment
# Documented defaults — override in /etc/qmanager/environment
Environment=PING_TARGET_1=http://www.gstatic.com/generate_204
Environment=PING_TARGET_2=http://cp.cloudflare.com/
Environment=PING_INTERVAL=5
Environment=FAIL_THRESHOLD=3
Environment=RECOVER_THRESHOLD=2
Environment=HISTORY_SIZE=60
Environment=CARRIER_FILE=/sys/class/net/rmnet_data0/carrier
TimeoutStopSec=10
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=3600
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
