# scripts/etc/systemd/system/qmanager-poller.service
# /lib/systemd/system/qmanager-poller.service
[Unit]
Description=QManager Modem Data Poller
After=qmanager-setup.service qmanager-ping.service
Wants=qmanager-ping.service

[Service]
Type=simple
ExecStartPre=/bin/sh -c 'i=0; while [ $i -lt 30 ] && [ ! -e /dev/smd11 ]; do sleep 1; i=$((i+1)); done; [ -e /dev/smd11 ] || { echo "AT device /dev/smd11 not found after 30s"; exit 1; }'
ExecStart=/usr/bin/qmanager_poller
EnvironmentFile=-/etc/qmanager/environment
TimeoutStartSec=35
TimeoutStopSec=10
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=3600
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
