# scripts/etc/udev/rules.d/99-qmanager-smd11.rules
# QManager: set /dev/smd11 to root:dialout mode 660 on every add event.
# Numeric prefix "99-" ensures this runs after OEM/vendor data_udev_rules.rules
# so we override their settings. ACTION=="add" prevents redundant firing on
# "change" or "remove" events.
#
# Note: SUBSYSTEM intentionally omitted. The subsystem string for smd11
# (rpmsg / qcom_smd / other) varies across Qualcomm kernels and has not been
# confirmed for both RM520N-GL and PRAIRE platforms. KERNEL=="smd11" is
# highly specific (Qualcomm SMD naming convention) and the helper is
# defensive ([ -c ] guard + exit 0), so false positives are harmless.
# To add SUBSYSTEM later: run `udevadm info /dev/smd11 | grep SUBSYSTEM` on
# the target device and prepend the match to the rule line.
KERNEL=="smd11", ACTION=="add", RUN+="/usr/lib/qmanager/qmanager_smd11_udev.sh"
